# .config/systemd/user/mail.service
[Unit]
Description=Mail service

[Service]
Type=simple
ExecStart=/bin/bash -c '/bin/pgrep geary || /usr/bin/geary --hidden'
