# .config/systemd/user/mega-sync.service
[Unit]
Description=MEGA sync service

[Service]
Type=simple
ExecStart=/bin/bash -c '/bin/pgrep mega-cmd-server || /usr/bin/mega-cmd-server'
