# .config/systemd/user/mega-sync.timer
[Unit]
Description=MEGA sync timer
 
[Timer]
Unit=mega-sync.service
OnCalendar=*-*-* *:*:00
Persistent=true
 
[Install]
WantedBy=timers.target
