# .config/systemd/user/mega-sync.service
[Unit]
Description=MEGA sync service

[Service]
Type=simple
ExecStart=bash -c 'pgrep mega-cmd-server || mega-cmd-server'
