# .config/systemd/user/mega-sync.service
[Unit]
Description=MEGA sync service

[Service]
Type=simple
ExecStart=/usr/bin/mega-cmd-server
