# .config/systemd/user/paperboy.timer
[Unit]
Description=PaperBoy timer
 
[Timer]
Unit=paperboy.service
OnCalendar=*-*-* *:*:0/5
Persistent=true
 
[Install]
WantedBy=timers.target
