# .config/systemd/user/paperboy.service
[Unit]
Description=PaperBoy service
ConditionPathExists=%h/.config/paperboy
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=paperboy -c %h/.config/paperboy
Restart=always
RestartSec=10

[Install]
WantedBy=xsession.target
