# .config/systemd/user/paperboy.service
[Unit]
Description=PaperBoy service

[Service]
Type=simple
ExecStart=%h/.local/bin/paperboy -c  %h/.config/paperboy
