# .config/systemd/user/paperboy.service
[Unit]
Description=PaperBoy service

[Service]
Type=simple
ExecStart=paperboy -c  %h/.config/paperboy
