# .config/systemd/user/taskwarrior-syncer.service
[Unit]
Description=TaskWarrior cloud syncer service

[Service]
Type=simple
ExecStart=timeout -k1 2 task synchronize
