# .config/systemd/user/taskwarrior-syncer.service
[Unit]
Description=TaskWarrior cloud syncer service

[Service]
Type=simple
ExecStart=/usr/bin/task synchronize
