# .config/systemd/user/taskwarrior-syncer.timer
[Unit]
Description=TaskWarrior cloud syncer timer
 
[Timer]
Unit=taskwarrior-syncer.service
OnCalendar=*-*-* *:*:00
Persistent=true
 
[Install]
WantedBy=timers.target
