# .config/systemd/user/wal-flusher.service
[Unit]
Description=Wal configuration files flusher service

[Service]
Type=simple
ExecStart=/bin/bash %h/.config/systemd/user/wal-flusher.sh
