# .config/systemd/user/xsession.target
[Unit]
Description=Xsession up
BindsTo=graphical-session.target
