# .config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Name=Conky
Description=Conky monitor
Exec=conky
